# .config/systemd/user/random-wallpaper.service
[Unit]
Description=Random Wallpaper Changer
[Service]
ExecStart=/usr/local/bin/random_wallpaper
Type=oneshot
